# recipes/wip/linux/linux-base/etc/inittab
ttyS0::respawn:/bin/login -f root
tty1::respawn:/bin/login -f root
::sysinit:/bin/hostname HOSTNAME
::sysinit:/bin/mount -a
